# etc/hello.conf
[program:hello]
command = gunicorn -b 0.0.0.0:8080 hello:app
directory = /home/dzen90/web/etc
user = www-data
